# systemd/rag-worker.service
[Unit]
Description=Immigration RAG — SQS Worker
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user/immigration-rag-claud-code-folder
EnvironmentFile=/home/ec2-user/immigration-rag-claud-code-folder/.env
ExecStart=/home/ec2-user/immigration-rag-claud-code-folder/.venv/bin/python worker.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
