# systemd/rag-worker@.service
[Unit]
Description=Immigration RAG — SQS Worker (%i)
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/immigration-rag
EnvironmentFile=/home/ubuntu/immigration-rag/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/immigration-rag/.venv/bin/python worker.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
