# systemd/rag-api.service
[Unit]
Description=Immigration RAG — Flask API
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user/immigration-rag-claud-code-folder
EnvironmentFile=/home/ec2-user/immigration-rag-claud-code-folder/.env
ExecStart=/home/ec2-user/immigration-rag-claud-code-folder/.venv/bin/gunicorn app:app --bind 0.0.0.0:5000 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
