# systemd/rag-api.service
[Unit]
Description=Immigration RAG — Flask API
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/immigration-rag
EnvironmentFile=/home/ubuntu/immigration-rag/.env
# Single worker: in-memory _s3_cache / _cache in app.py are per process — 2 workers => two caches (stale / divergent S3 browse).
ExecStart=/home/ubuntu/immigration-rag/.venv/bin/gunicorn app:app --bind 0.0.0.0:5000 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
